# tcpdump.service
[Unit]
Description= Tcpdump runner

[Service]
ExecStart=/usr/local/bin/tcpdump.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
