# tcpdump.service
[Unit]
Description= Tcpdump runner

[Service]
ExecStart=/home/ubuntu/tcpdump.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
